# build/deb/lib/systemd/system/deadendproxy.service
[Unit]
Description=DeadEndProxy Service
After=network.target
[Service]
ExecStart=/usr/local/bin/deadendproxy -port-http 80 -port-proxy 443 -config /etc/deadendproxy/config.yaml
Restart=always
User=root
WorkingDirectory=/etc/deadendproxy
[Install]
WantedBy=multi-user.target
